# debian/cron.daily
#!/bin/sh

SXID_OPTS=

if [ -x /usr/bin/sxid ]; then
	/usr/bin/sxid ${SXID_OPTS}
fi
