# docs/sxid.cron.example
# Sample crontab entry to run every day at 4am
0 4 * * * /usr/bin/sxid
